# web/packages/services/dev/jn-mirror-socket.service
[Unit]
Description=Mirror Websocket
After=redis.service
Requires=redis.service

[Service]
TimeoutStartSec=0
ExecStart=/var/www/Sockets/bin/start.sh

[Install]
WantedBy=multi-user.target

