# web/packages/services/dev/jn-mirror-sockets.target
[Unit]
Description=Mirror Websocket
After=jn-mirror-socket@0.service jn-mirror-socket@1.service jn-mirror-socket@2.service
Requires=jn-mirror-socket@0.service jn-mirror-socket@1.service jn-mirror-socket@2.service

[Install]
WantedBy=multi-user.target
