# web/Backend/config/backend.cron
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=root
HOME=/

0 16 * * * root php /var/www/jukebox/Backend/writer.php InitialVevoArtistsImport
0 20 * * * root php /var/www/jukebox/Backend/writer.php InitialVevoArtistsVideosImport
0 0 * * * root php /var/www/jukebox/Backend/writer.php Initial
